# config/systemctl/massbit-gateway/update-worker-ip.timer
[Unit]
Description=Run update-worker-ip.service

[Timer]
OnBootSec=30
OnUnitActiveSec=5

[Install]
WantedBy=timers.target